# roles/prisma_backup/templates/prisma-backup.timer.j2
# Managed by Ansible — prisma-twistlock-install / prisma_backup
[Unit]
Description=Run Prisma Console backup nightly

[Timer]
OnCalendar={{ prisma_backup_schedule }}
Persistent=true
RandomizedDelaySec={{ prisma_backup_jitter_sec }}

[Install]
WantedBy=timers.target
